# scripts/caddypanel.service
[Unit]
Description=CaddyPanel - Caddy Reverse Proxy Management Panel
Documentation=https://github.com/caddypanel/caddypanel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=caddypanel
Group=caddypanel
ExecStart=/usr/local/bin/caddypanel
WorkingDirectory=/var/lib/caddypanel
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Environment
EnvironmentFile=-/etc/caddypanel/caddypanel.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/caddypanel /var/log/caddypanel
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=caddypanel

[Install]
WantedBy=multi-user.target
